# rootfs/etc/services.d/usbip/run
#!/usr/bin/with-contenv bashio
bashio::log.info "usbip service: running simple-usbip-client.sh"
exec /usr/local/bin/simple-usbip-client.sh
